# .config/systemd/user/dotfiles-sync.service
[Unit]
Description=Apply dotfiles managed state

[Service]
Type=oneshot
ExecStart=%h/ghq/github.com/schroneko/dotfiles/scripts/dotfiles-sync.sh --apply-only
